# .config/systemd/user/notify.service
[Unit]
Description=Wasser trinken

[Service]
Type=oneshot
ExecStart=/home/saj/.config/systemd/user/wasser.sh

[Install]
WantedBy=multi-user.target
